# examples/echo/systemd/echo.target
[Unit]
Description=Erlang Echo service
Requires=echo-tcp.socket echo-udp.socket
After=echo-tcp.socket echo-udp.socket
