# examples/plug/systemd/plug.socket
[Unit]
Description=Listening socket
Requires=sockets.target

[Socket]
ListenStream=80
BindIPv6Only=both
ReusePort=true
NoDelay=true
