# examples/plug/systemd/plug.target
[Unit]
Description=Simple HTTP server example in Elixir
Requires=plug.socket
After=plug.socket
